# content/service/rclone_webdav/run
#!/bin/sh

DRIVENAME="$(grep ^drive-name /mnt/data/config/script.conf | cut -d= -f2-)"

# Run rclone
htpasswd -csb /.aria2allinoneworkdir/htpasswd_webdav ${GLOBAL_USER} ${GLOBAL_PASSWORD}
exec 2>&1
exec rclone serve webdav ${DRIVENAME}:/ \
--htpasswd=/.aria2allinoneworkdir/htpasswd_webdav \
--addr=localhost:61805 \
--buffer-size 4M \
--config /mnt/data/config/rclone.conf \
--use-mmap \
--vfs-cache-mode writes 