# content/service/metube/run
#!/bin/sh

cd /app
exec 2>&1
exec python3 /app/app/main.py