# content/service/ttyd/run
#!/bin/sh

exec 2>&1
exec ttyd -i 127.0.0.1 -P 3 -t fontSize=17 -t titleFixed=Heroku_ttyd -t 'theme={"background": "black"}' bash --rcfile ./.bashrc