# content/service/backup_files/run
#!/bin/sh

sleep 300

DRIVE_NAME_AUTO="$(sed -n '1p' /mnt/data/config/rclone.conf | sed "s/\[//g" | sed "s/\]//g")"
if [ "${DRIVE_NAME}" = "auto" ]; then
    DRIVENAME=${DRIVE_NAME_AUTO}
else 
    DRIVENAME=${DRIVE_NAME}
fi

exec 2>&1
rclone rc --user "${USER}" --pass "${PASSWORD}" --rc-addr=localhost:56802 operations/movefile srcFs=/.aria2allinoneworkdir/aria2 srcRemote=aria2.session dstFs="${DRIVENAME}":/herokuallinoneapp dstRemote=aria2.session _async=true
rclone rc --user "${USER}" --pass "${PASSWORD}" --rc-addr=localhost:56802 operations/movefile srcFs=/.aria2allinoneworkdir/aria2 srcRemote=dht.dat dstFs="${DRIVENAME}":/herokuallinoneapp dstRemote=dht.dat _async=true
rclone rc --user "${USER}" --pass "${PASSWORD}" --rc-addr=localhost:56802 operations/movefile srcFs=/.aria2allinoneworkdir/aria2 srcRemote=aria2.conf dstFs="${DRIVENAME}":/herokuallinoneapp dstRemote=aria2.conf _async=true
rclone rc --user "${USER}" --pass "${PASSWORD}" --rc-addr=localhost:56802 operations/movefile srcFs=/.aria2allinoneworkdir/aria2 srcRemote=script.conf dstFs="${DRIVENAME}":/herokuallinoneapp dstRemote=script.conf _async=true
rclone rc --user "${USER}" --pass "${PASSWORD}" --rc-addr=localhost:56802 operations/movefile srcFs=/.aria2allinoneworkdir srcRemote=filebrowser.db dstFs="${DRIVENAME}":/herokuallinoneapp dstRemote=filebrowser.db _async=true